# targets/imx6ullevk/rootfs_overlay/etc/systemd/system/tee-supplicant.service
[Unit]
Description=Tee Supplicant

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=/usr/sbin/tee-supplicant
User=root

[Install]
WantedBy=multi-user.target
